# deploy/linux/did-intel.service
[Unit]
Description=DID Intel API Server
Documentation=https://github.com/siamakanda/DIDRepChecker
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/did-intel
Environment="PATH=/opt/did-intel/venv/bin:/usr/bin"
Environment="DIDINTEL_API_HOST=127.0.0.1"
Environment="DIDINTEL_API_PORT=8000"
ExecStart=/opt/did-intel/venv/bin/uvicorn did_intel.api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
